# root/etc/services.d/unitd/run
#!/command/with-contenv bash

echo "starting unitd"
exec /usr/sbin/unitd --no-daemon --log /dev/stdout;


